# encfs/s6-overlay/s6-rc.d/storage/run
#!/command/with-contenv bash
exec /app/run.sh
#exec s6-setuidgid securefs /app/run.sh
